# RASPBERRY-PI/HardwareTrigger/rpi3-camera-prototype.conf
[program:rpi3-camera-prototype]
directory=/home/pi/Programming
command=/usr/bin/python ./camtrigger.py 
autostart=true
autorestart=true
stopsignal=INT
stopasgroup=true
killasgroup=true
user=pi
